# supervisord.conf
[supervisord]
nodaemon=true

[program:node]
autorestart=true
command=/usr/bin/node server --port=9001
